# crontab
# Run all scrapers every 4 hours (same as container startup)
0 */4 * * * docker exec cryptorank-scraper sh -c "npm run funding-rounds && npm run project-details && npm run generate-test-data && npm run tier1 && npm run tier2" >> /proc/1/fd/1 2>&1
